# terminal-gateway.service
[Unit]
Description=Terminal Gateway Server
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/home/webterminal
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=terminal-gateway

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target